# rc.local
#!/bin/rc

sysctl -p /etc/sysctl.conf

hostname oasis &
{
	ip link set lo up
	# add other networking stuff here
} &

wait
