# rc.local
#!/bin/rc

sysctl -p /etc/sysctl.conf

hostname oasis &
ip link set lo up &
ip link set eth0 up &

wait
